# multiwan-nft/files/etc/hotplug.d/iface/15-multiwan-nft
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/multiwan-nft/multiwan_nft.sh

initscript=/etc/init.d/multiwan-nft
. /lib/functions/procd.sh


SCRIPTNAME="multiwan_nft-hotplug"

# Enhanced debug logging for troubleshooting
LOG debug "multiwan_nft hotplug started: ACTION=$ACTION INTERFACE=$INTERFACE DEVICE=$DEVICE MULTIWAN_NFT_STARTUP=$MULTIWAN_NFT_STARTUP FIRSTCONNECT=$FIRSTCONNECT MULTIWAN_NFT_SHUTDOWN=$MULTIWAN_NFT_SHUTDOWN"

[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifdown" ] || [ "$ACTION" = "connected" ] || [ "$ACTION" = "disconnected" ] || exit 1
[ -n "$INTERFACE" ] || exit 2
[ "$FIRSTCONNECT" = "1" ] || [ "$MULTIWAN_NFT_SHUTDOWN" = "1" ] && exit 0

if { [ "$ACTION" = "ifup" ] || [ "$ACTION" = "connected" ] ; } && [ -z "$DEVICE" ]; then
	LOG notice "$ACTION called on $INTERFACE with no device set"
	exit 3
fi

[ "$MULTIWAN_NFT_STARTUP" = "init" ] || procd_lock

multiwan_nft_init

/etc/init.d/multiwan-nft running || {
	LOG notice "multiwan_nft hotplug $ACTION on $INTERFACE not called because globally disabled"
	exit 0
}

$NFT list chain $NFT_FAMILY $NFT_TABLE multiwan_nft_hook >/dev/null 2>&1 || {
	LOG warn "hotplug called on $INTERFACE before multiwan_nft has been set up"
	exit 0
}

if [ "$MULTIWAN_NFT_STARTUP" != "init" ] && [ "$ACTION" = "ifup" ]; then
	multiwan_nft_set_user_iface_rules_nft $INTERFACE $DEVICE
fi

config_get_bool enabled $INTERFACE 'enabled' '0'
[ "${enabled}" -eq 1 ] || {
	LOG notice "multiwan_nft hotplug on $INTERFACE not called because interface disabled"
	exit 0
}

config_get initial_state $INTERFACE initial_state "online"
if [ "$initial_state" = "offline" ]; then
	status=""
	if [ -r "$MULTIWAN_NFT_TRACK_STATUS_DIR/$INTERFACE/STATUS" ]; then
		{ IFS= read -r status < "$MULTIWAN_NFT_TRACK_STATUS_DIR/$INTERFACE/STATUS" || [ -n "$status" ]; } 2>/dev/null || status=
	fi
	status="${status:-unknown}"
	[ "$status" = "online" ] || status=offline
else
	status=online
fi

LOG notice "Execute $ACTION event on interface $INTERFACE (${DEVICE:-unknown})"

status_dir="$MULTIWAN_NFT_TRACK_STATUS_DIR/$INTERFACE"
mkdir -p "$status_dir"
old_state="$(multiwan_nft_get_iface_hotplug_state "$INTERFACE")"
previous_session_action=""
if [ -r "$status_dir/SESSION_ACTION" ]; then
	{ IFS= read -r previous_session_action < "$status_dir/SESSION_ACTION" || [ -n "$previous_session_action" ]; } 2>/dev/null || previous_session_action=
fi
case "$previous_session_action" in
	purged|already-purged) failure_already_purged=1 ;;
	*) failure_already_purged=0 ;;
esac
record_transition() {
	printf '%s\n' "$2" > "$status_dir/$1"
}
record_last_failure() {
	printf '%s\n' "$2" > "$status_dir/LAST_$1"
}
capture_qos_snapshot() {
	local snapshot_file="$status_dir/LAST_QOS_SNAPSHOT"
	local qos_dev snapshot_devices=0

	if ! command -v tc >/dev/null 2>&1 || [ -z "$DEVICE" ]; then
		printf '%s\n' "tc or interface unavailable" > "$snapshot_file"
		record_last_failure QOS_BACKLOG "unknown"
		return 0
	fi

	{
		date
		for qos_dev in "$DEVICE" "ifb-$DEVICE"; do
			ip link show dev "$qos_dev" >/dev/null 2>&1 || continue
			snapshot_devices=$((snapshot_devices + 1))
			printf '%s\n' "=== $qos_dev ==="
			ip -s link show dev "$qos_dev" 2>/dev/null
			tc -s qdisc show dev "$qos_dev" 2>/dev/null
			tc -s class show dev "$qos_dev" 2>/dev/null | awk '
				/^class hfsc 1:(11|13|15) / { selected = 1; print; next }
				/^class / { selected = 0 }
				selected { print }
			'
		done
	} > "$snapshot_file"

	if [ "$snapshot_devices" -eq 0 ]; then
		record_last_failure QOS_BACKLOG "unknown"
	elif grep -Eq 'backlog ([1-9][0-9]*[bKMG]|[0-9]+b [1-9][0-9]*p)' "$snapshot_file"; then
		record_last_failure QOS_BACKLOG "present"
	else
		record_last_failure QOS_BACKLOG "none"
	fi
}
record_failure_context() {
	record_last_failure TRANSITION "$ACTION"
	record_last_failure TRANSITION_AT "$(date +%s)"
	record_last_failure OFFLINE_REASON "${OFFLINE_REASON:-$ACTION}"
	capture_qos_snapshot
}
run_recovery_cleanup() {
	local current_state
	current_state="$(multiwan_nft_get_iface_hotplug_state "$INTERFACE")"
	if [ "$old_state" != "offline" ] || [ "$current_state" != "online" ] || [ "$MULTIWAN_NFT_STARTUP" = "init" ]; then
		record_transition SESSION_ACTION "not-required"
		record_transition SESSION_COUNT "0"
		return 0
	fi

	if multiwan_nft_flush_conntrack "$INTERFACE" "$ACTION"; then
		if [ "${MULTIWAN_NFT_FLUSH_MATCHED:-0}" -eq 1 ]; then
			record_transition SESSION_ACTION "purged"
			record_transition SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
		else
			record_transition SESSION_ACTION "not-required"
			record_transition SESSION_COUNT "0"
		fi
	else
		record_transition SESSION_ACTION "purge-failed"
		return 1
	fi
}
case "$ACTION" in
	connected|disconnected)
		record_transition PROBE_TRANSITION "$ACTION"
		record_transition PROBE_TRANSITION_AT "$(date +%s)"
		[ "$ACTION" = "disconnected" ] && record_transition OFFLINE_REASON "${OFFLINE_REASON:-probe}" || record_transition OFFLINE_REASON "none"
		;;
	ifdown) record_transition OFFLINE_REASON "ifdown" ;;
	*) record_transition OFFLINE_REASON "none" ;;
esac

case "$ACTION" in
	connected)
		LOG debug "Processing connected: committing policy before restoring online lookup rules"
		multiwan_nft_set_iface_hotplug_state $INTERFACE "online"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			LOG warn "connected event on $INTERFACE left the previous policy active"
			exit 1
		}
		multiwan_nft_create_iface_rules $INTERFACE $DEVICE || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			multiwan_nft_set_policies_nft >/dev/null 2>&1
			[ "$old_state" = "offline" ] && multiwan_nft_set_tracker_oif_rules "$INTERFACE" "$DEVICE" >/dev/null 2>&1
			record_transition POLICY_RESULT "rolled-back-rule-failure"
			LOG warn "connected event on $INTERFACE rolled policy back after online lookup rules failed"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		record_transition SESSION_ACTION "not-required"
		record_transition SESSION_COUNT "0"
		record_transition STICKY_RESULT "not-required"
		;;
	ifup)
		LOG debug "Processing ifup: creating iface rules and nft chains"
		multiwan_nft_create_iface_nft $INTERFACE $DEVICE || {
			LOG warn "ifup event on $INTERFACE failed while rebuilding nft interface state"
			exit 1
		}
		multiwan_nft_create_iface_rules $INTERFACE $DEVICE || {
			LOG warn "ifup event on $INTERFACE failed while rebuilding ip rules"
			exit 1
		}
		multiwan_nft_set_iface_hotplug_state $INTERFACE "$status"
		if [ "$MULTIWAN_NFT_STARTUP" != "init" ]; then
			multiwan_nft_create_iface_route $INTERFACE $DEVICE || {
				LOG warn "ifup event on $INTERFACE skipped dependent rebuilds because route table update failed"
				exit 1
			}
			multiwan_nft_set_general_rules
			# Repopulate connected sets to ensure local traffic bypasses multiwan_nft
			multiwan_nft_set_connected_ipv4
			[ "$NO_IPV6" -eq 0 ] && multiwan_nft_set_connected_ipv6
			if [ "$status" = "online" ]; then
				multiwan_nft_set_policies_nft || {
					multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
					record_transition POLICY_RESULT "failed"
					LOG warn "ifup event on $INTERFACE left the previous policy active"
					exit 1
				}
				record_transition POLICY_RESULT "applied"
			fi
		fi
		[ "$ACTION" = ifup ] && procd_running multiwan_nft "track_$INTERFACE" && procd_send_signal multiwan_nft "track_$INTERFACE" USR2
		;;
	disconnected)
		LOG debug "Processing disconnected: committing policy before scoped cleanup"
		[ "$old_state" != "offline" ] && record_failure_context
		multiwan_nft_set_iface_hotplug_state $INTERFACE "offline"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			record_transition SESSION_ACTION "preserved-policy-failure"
			record_last_failure POLICY_RESULT "failed"
			record_last_failure SESSION_ACTION "preserved-policy-failure"
			LOG warn "disconnected event on $INTERFACE retained policy and conntrack state"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		record_last_failure POLICY_RESULT "applied"
		if multiwan_nft_set_tracker_oif_rules $INTERFACE $DEVICE; then
			record_transition RECOVERY_RULE_RESULT "tracker-only"
		else
			record_transition RECOVERY_RULE_RESULT "failed"
			LOG warn "disconnected event on $INTERFACE blocked its mark but could not install every tracker-only recovery rule"
		fi
		if [ "$old_state" != "offline" ] || [ "$failure_already_purged" -ne 1 ]; then
			if multiwan_nft_delete_iface_set_entries $INTERFACE; then
				record_transition STICKY_RESULT "removed"
				record_transition STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
				record_last_failure STICKY_RESULT "removed"
				record_last_failure STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
			else
				record_transition STICKY_RESULT "failed"
				record_transition STICKY_REMOVED "failed"
				record_last_failure STICKY_RESULT "failed"
				record_last_failure STICKY_REMOVED "failed"
				LOG warn "disconnected event on $INTERFACE could not remove all sticky entries; the offline mark remains unreachable"
			fi
			if multiwan_nft_flush_conntrack_iface $INTERFACE; then
				record_transition SESSION_ACTION "purged"
				record_transition SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
				record_last_failure SESSION_ACTION "purged"
				record_last_failure SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
			else
				record_transition SESSION_ACTION "purge-failed"
				record_last_failure SESSION_ACTION "purge-failed"
				exit 1
			fi
		else
			record_transition SESSION_ACTION "already-purged"
		fi
		;;
	ifdown)
		LOG debug "Processing ifdown: committing policy before physical cleanup"
		[ "$old_state" != "offline" ] && record_failure_context
		multiwan_nft_set_iface_hotplug_state $INTERFACE "offline"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			record_transition SESSION_ACTION "preserved-policy-failure"
			record_last_failure POLICY_RESULT "failed"
			record_last_failure SESSION_ACTION "preserved-policy-failure"
			LOG warn "ifdown event on $INTERFACE retained previous MultiWAN state"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		record_last_failure POLICY_RESULT "applied"
		if [ "$old_state" != "offline" ] || [ "$failure_already_purged" -ne 1 ]; then
			if multiwan_nft_delete_iface_set_entries $INTERFACE; then
				record_transition STICKY_RESULT "removed"
				record_transition STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
				record_last_failure STICKY_RESULT "removed"
				record_last_failure STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
			else
				record_transition STICKY_RESULT "failed"
				record_transition STICKY_REMOVED "failed"
				record_last_failure STICKY_RESULT "failed"
				record_last_failure STICKY_REMOVED "failed"
			fi
			if multiwan_nft_flush_conntrack_iface $INTERFACE; then
				record_transition SESSION_ACTION "purged"
				record_transition SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
				record_last_failure SESSION_ACTION "purged"
				record_last_failure SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
			else
				record_transition SESSION_ACTION "purge-failed"
				record_last_failure SESSION_ACTION "purge-failed"
			fi
		else
			record_transition SESSION_ACTION "already-purged"
		fi
		multiwan_nft_delete_iface_rules $INTERFACE
		multiwan_nft_delete_iface_route $INTERFACE
		multiwan_nft_delete_iface_nft $INTERFACE
		procd_running multiwan_nft "track_$INTERFACE" && procd_send_signal multiwan_nft "track_$INTERFACE" USR1
	;;
esac

case "$ACTION" in
	connected|ifup) run_recovery_cleanup || exit 1 ;;
esac

LOG debug "multiwan_nft hotplug completed: ACTION=$ACTION INTERFACE=$INTERFACE"

exit 0
