# multiwan-nft/files/etc/hotplug.d/iface/15-multiwan-nft
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh
. /lib/multiwan-nft/multiwan_nft.sh

initscript=/etc/init.d/multiwan-nft
. /lib/functions/procd.sh


SCRIPTNAME="multiwan_nft-hotplug"

# Enhanced debug logging for troubleshooting
LOG debug "multiwan_nft hotplug started: ACTION=$ACTION INTERFACE=$INTERFACE DEVICE=$DEVICE MULTIWAN_NFT_STARTUP=$MULTIWAN_NFT_STARTUP FIRSTCONNECT=$FIRSTCONNECT MULTIWAN_NFT_SHUTDOWN=$MULTIWAN_NFT_SHUTDOWN"

[ "$ACTION" = "ifup" ] || [ "$ACTION" = "ifdown" ] || [ "$ACTION" = "connected" ] || [ "$ACTION" = "disconnected" ] || exit 1
[ -n "$INTERFACE" ] || exit 2
[ "$FIRSTCONNECT" = "1" ] || [ "$MULTIWAN_NFT_SHUTDOWN" = "1" ] && exit 0

if { [ "$ACTION" = "ifup" ] || [ "$ACTION" = "connected" ] ; } && [ -z "$DEVICE" ]; then
	LOG notice "$ACTION called on $INTERFACE with no device set"
	exit 3
fi

[ "$MULTIWAN_NFT_STARTUP" = "init" ] || procd_lock

multiwan_nft_init

/etc/init.d/multiwan-nft running || {
	LOG notice "multiwan_nft hotplug $ACTION on $INTERFACE not called because globally disabled"
	exit 0
}

$NFT list chain $NFT_FAMILY $NFT_TABLE multiwan_nft_hook &>/dev/null || {
	LOG warn "hotplug called on $INTERFACE before multiwan_nft has been set up"
	exit 0
}

if [ "$MULTIWAN_NFT_STARTUP" != "init" ] && [ "$ACTION" = "ifup" ]; then
	multiwan_nft_set_user_iface_rules_nft $INTERFACE $DEVICE
fi

config_get_bool enabled $INTERFACE 'enabled' '0'
[ "${enabled}" -eq 1 ] || {
	LOG notice "multiwan_nft hotplug on $INTERFACE not called because interface disabled"
	exit 0
}

config_get initial_state $INTERFACE initial_state "online"
if [ "$initial_state" = "offline" ]; then
	status=$(cat $MULTIWAN_NFT_TRACK_STATUS_DIR/$INTERFACE/STATUS 2>/dev/null || echo unknown)
	[ "$status" = "online" ] || status=offline
else
	status=online
fi

LOG notice "Execute $ACTION event on interface $INTERFACE (${DEVICE:-unknown})"

status_dir="$MULTIWAN_NFT_TRACK_STATUS_DIR/$INTERFACE"
mkdir -p "$status_dir"
old_state="$(multiwan_nft_get_iface_hotplug_state "$INTERFACE")"
record_transition() {
	printf '%s\n' "$2" > "$status_dir/$1"
}
case "$ACTION" in
	disconnected) record_transition OFFLINE_REASON "${OFFLINE_REASON:-probe}" ;;
	ifdown) record_transition OFFLINE_REASON "ifdown" ;;
	*) record_transition OFFLINE_REASON "none" ;;
esac

case "$ACTION" in
	connected)
		LOG debug "Processing connected: committing policy before restoring online lookup rules"
		multiwan_nft_set_iface_hotplug_state $INTERFACE "online"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			LOG warn "connected event on $INTERFACE left the previous policy active"
			exit 1
		}
		multiwan_nft_create_iface_rules $INTERFACE $DEVICE || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			multiwan_nft_set_policies_nft >/dev/null 2>&1
			[ "$old_state" = "offline" ] && multiwan_nft_set_tracker_oif_rules "$INTERFACE" "$DEVICE" >/dev/null 2>&1
			record_transition POLICY_RESULT "rolled-back-rule-failure"
			LOG warn "connected event on $INTERFACE rolled policy back after online lookup rules failed"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		record_transition SESSION_ACTION "not-required"
		record_transition SESSION_COUNT "0"
		record_transition STICKY_RESULT "not-required"
		;;
	ifup)
		LOG debug "Processing ifup: creating iface rules and nft chains"
		multiwan_nft_create_iface_nft $INTERFACE $DEVICE || {
			LOG warn "ifup event on $INTERFACE failed while rebuilding nft interface state"
			exit 1
		}
		multiwan_nft_create_iface_rules $INTERFACE $DEVICE || {
			LOG warn "ifup event on $INTERFACE failed while rebuilding ip rules"
			exit 1
		}
		multiwan_nft_set_iface_hotplug_state $INTERFACE "$status"
		if [ "$MULTIWAN_NFT_STARTUP" != "init" ]; then
			multiwan_nft_create_iface_route $INTERFACE $DEVICE || {
				LOG warn "ifup event on $INTERFACE skipped dependent rebuilds because route table update failed"
				exit 1
			}
			multiwan_nft_set_general_rules
			# Repopulate connected sets to ensure local traffic bypasses multiwan_nft
			multiwan_nft_set_connected_ipv4
			[ "$NO_IPV6" -eq 0 ] && multiwan_nft_set_connected_ipv6
			if [ "$status" = "online" ]; then
				multiwan_nft_set_policies_nft || {
					multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
					record_transition POLICY_RESULT "failed"
					LOG warn "ifup event on $INTERFACE left the previous policy active"
					exit 1
				}
				record_transition POLICY_RESULT "applied"
			fi
		fi
		[ "$ACTION" = ifup ] && procd_running multiwan_nft "track_$INTERFACE" && procd_send_signal multiwan_nft "track_$INTERFACE" USR2
		;;
	disconnected)
		LOG debug "Processing disconnected: committing policy before scoped cleanup"
		multiwan_nft_set_iface_hotplug_state $INTERFACE "offline"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			record_transition SESSION_ACTION "preserved-policy-failure"
			LOG warn "disconnected event on $INTERFACE retained policy and conntrack state"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		if multiwan_nft_set_tracker_oif_rules $INTERFACE $DEVICE; then
			record_transition RECOVERY_RULE_RESULT "tracker-only"
		else
			record_transition RECOVERY_RULE_RESULT "failed"
			LOG warn "disconnected event on $INTERFACE blocked its mark but could not install every tracker-only recovery rule"
		fi
		if multiwan_nft_delete_iface_set_entries $INTERFACE; then
			record_transition STICKY_RESULT "removed"
			record_transition STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
		else
			record_transition STICKY_RESULT "failed"
			record_transition STICKY_REMOVED "failed"
			LOG warn "disconnected event on $INTERFACE could not remove all sticky entries; the offline mark remains unreachable"
		fi
		if multiwan_nft_flush_conntrack_iface $INTERFACE; then
			record_transition SESSION_ACTION "purged"
			record_transition SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
		else
			record_transition SESSION_ACTION "purge-failed"
			exit 1
		fi
		;;
	ifdown)
		LOG debug "Processing ifdown: committing policy before physical cleanup"
		multiwan_nft_set_iface_hotplug_state $INTERFACE "offline"
		multiwan_nft_set_policies_nft || {
			multiwan_nft_set_iface_hotplug_state "$INTERFACE" "$old_state"
			record_transition POLICY_RESULT "failed"
			record_transition SESSION_ACTION "preserved-policy-failure"
			LOG warn "ifdown event on $INTERFACE retained previous MultiWAN state"
			exit 1
		}
		record_transition POLICY_RESULT "applied"
		previous_session_action="$(cat "$status_dir/SESSION_ACTION" 2>/dev/null)"
		if [ "$old_state" != "offline" ] || [ "$previous_session_action" != "purged" ]; then
			if multiwan_nft_delete_iface_set_entries $INTERFACE; then
				record_transition STICKY_RESULT "removed"
				record_transition STICKY_REMOVED "${MULTIWAN_NFT_STICKY_DELETED:-0}"
			else
				record_transition STICKY_RESULT "failed"
				record_transition STICKY_REMOVED "failed"
			fi
			if multiwan_nft_flush_conntrack_iface $INTERFACE; then
				record_transition SESSION_ACTION "purged"
				record_transition SESSION_COUNT "${MULTIWAN_NFT_CONNTRACK_PURGED:-0}"
			else
				record_transition SESSION_ACTION "purge-failed"
			fi
		else
			record_transition SESSION_ACTION "already-purged"
		fi
		multiwan_nft_delete_iface_rules $INTERFACE
		multiwan_nft_delete_iface_route $INTERFACE
		multiwan_nft_delete_iface_nft $INTERFACE
		procd_running multiwan_nft "track_$INTERFACE" && procd_send_signal multiwan_nft "track_$INTERFACE" USR1
	;;
esac

case "$ACTION" in
	connected|ifup) multiwan_nft_flush_conntrack "$INTERFACE" "$ACTION" ;;
esac

LOG debug "multiwan_nft hotplug completed: ACTION=$ACTION INTERFACE=$INTERFACE"

exit 0
